# scripts/com.j105.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.j105.backup.plist — launchd agent for daily Pi-to-Mac backup.
  Installed by scripts/setup-backup-mac.sh into ~/Library/LaunchAgents/.

  Runs backup.sh every day at 03:00 local time.
  Logs stdout+stderr to ~/backups/j105-logger/backup.log.

  Manual control (after installation):
    launchctl start  com.j105.backup   # run now
    launchctl stop   com.j105.backup   # cancel running job
    launchctl unload ~/Library/LaunchAgents/com.j105.backup.plist  # disable
    launchctl load   ~/Library/LaunchAgents/com.j105.backup.plist  # re-enable
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.j105.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>BACKUP_SCRIPT_PATH</string>
  </array>

  <!-- Daily at 03:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Log file — directory must exist before launchd writes to it -->
  <key>StandardOutPath</key>
  <string>BACKUP_LOG_PATH</string>
  <key>StandardErrorPath</key>
  <string>BACKUP_LOG_PATH</string>

  <!-- Run even if the scheduled time was missed (e.g. Mac was asleep) -->
  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>HOME_PATH</string>
  </dict>
</dict>
</plist>
